# desktop/etc/aldos/home.mount
[Mount]
Where=/home
What=tmpfs
Type=tmpfs
Options=rw,gid=users,size=256M

[Install]
RequiredBy=local-fs.target
